# scoreboard.service
[Unit]
Description=scoreboard
After=network-online.target
[Service]
ExecStart=/home/pi/scoreboard/app.py
[Install]
WantedBy=multi-user.target